# systemd/brightctl-day.service
# ~/.config/systemd/user/brightctl-day.service
[Unit]
Description=Run brightctl.sh with 100 at 08:00
After=default.target

[Service]
Type=oneshot
ExecStart=/home/%u/.local/bin/brightctl.sh 100
