# systemd/brightctl-night.service
# ~/.config/systemd/user/brightctl-night.service
[Unit]
Description=Decrease display brightness to 33 over 30 minutes using brightctl.sh
After=default.target

[Service]
Type=oneshot
ExecStart=/home/%u/.local/bin/brightctl.sh 33 1800
