# systemd/brightctl-day.service
# ~/.config/systemd/user/brightctl-day.service
[Unit]
Description=Run brightctl.sh to raise display brightness to 100 immediately
After=default.target

[Service]
Type=oneshot
ExecStart=/home/%u/.local/bin/brightctl.sh 100
