# systemd/brightctl-night.service
# ~/.config/systemd/user/brightctl-night.service
[Unit]
Description=Run brightctl.sh to lower display brightness to 33 in 30 minutes
After=default.target

[Service]
Type=oneshot
ExecStart=/home/%u/.local/bin/brightctl.sh 33 1800
