# systemd/brightctl-day.timer
# ~/.config/systemd/user/brightctl-day.timer
[Unit]
Description=Run brightctl-day.service at 08:00 to increase display brightness to 100 immediately 

[Timer]
OnCalendar=*-*-* 08:00:00
Persistent=false
RandomizedDelaySec=300
Unit=brightctl-day.service

[Install]
WantedBy=timers.target
