# systemd/brightctl-day.service
# ~/.config/systemd/user/brightctl-day.service
[Unit]
Description=Increase display brightness to 100 immediately using brightctl.sh
After=default.target

[Service]
Type=oneshot
ExecStart=/home/%u/.local/bin/brightctl.sh 100
